# safekeep.cron
#!/bin/sh

PATH=/sbin:/bin:/usr/sbin:/usr/bin

safekeep -v --server
