# safekeep.cron
#!/bin/sh

PATH=/sbin:/bin:/usr/sbin:/usr/bin

if ls /etc/safekeep/backup.d/*.backup 2>&1 1>/dev/null; then
    safekeep -v --server
fi
