# safekeep.cron
#!/bin/sh

PATH=/sbin:/bin:/usr/sbin:/usr/bin

sudo -H -u safekeep safekeep --server

